# agentguard-sensor/install/macos/com.agentguard.sensor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
<key>Label</key><string>com.agentguard.sensor</string>
<key>ProgramArguments</key><array><string>/usr/local/bin/agentguard-sensor</string><string>--config</string><string>/Library/Application Support/AgentGuard/config.yml</string><string>watch</string></array>
<key>RunAtLoad</key><true/>
<key>KeepAlive</key><true/>
<key>StandardOutPath</key><string>/Library/Logs/AgentGuard/stdout.log</string>
<key>StandardErrorPath</key><string>/Library/Logs/AgentGuard/stderr.log</string>
</dict></plist>
